# auth-iitk.conf
[program:auth-iitk]
directory=/opt/
command=python3 auth-iitk.py
autostart=true
autorestart=true
